# cron.d/log4k
00 23 * * * root /usr/local/bin/python /app/main.py >> /var/log/python.log
