# cron.d/log4k
00 23 * * * root /usr/local/bin/python /app/log4k/__main__.py >> /var/log/python.log
